# docker/files/etc/xinetd.d/cvspserver
# default: off
# description: The cvspserver server is a good addition to an ftp server, as it \
#	allows crc checksumming etc.
service cvspserver
{
	disable		= no
	socket_type     = stream
	wait            = no
	user            = root
	server          = /usr/bin/cvs
	server_args     = -f --allow-root=/cvs-repository --allow-root=/cvs-zopeorg pserver
	log_on_failure  += USERID
}
