# docker/files/etc/xinetd.d/svn
service svn
{
   socket_type     = stream
   protocol        = tcp
   wait            = no
   user            = nobody
   group           = zopesvn
   server          = /usr/local/bin/svnserve
   server_args     = -i
   disable         = no
} 
